# tests/integration/files/playbook/awsnycast.service.j2
[Unit]
Description=Job that runs AWSnycast

[Service]
Type=simple
ExecStart=/usr/local/bin/awsnycast